# templates/default/sand-agent.init.conf.erb
# Generated by Chef
#

description "SAND - SAND is Awesome Network Daemon"
author "leo@scalingo.com"

start on runlevel [2345]
stop on runlevel [016]

respawn

<% @env.each do |item| %>
env <%= item %>
<% end %>

exec <%= @target %>
